# core/console/systemd/corvin-operator-ui-watchdog.timer
[Unit]
Description=Corvin Operator UI — health probe every 60s
Documentation=https://github.com/anthropics/corvinos

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
AccuracySec=5s
Unit=corvin-operator-ui-watchdog.service

[Install]
WantedBy=timers.target
