# operator/bridges/slack/systemd/corvin-voice-bridge-slack.service
[Unit]
Description=corvin-voice Slack bridge daemon (Bolt Socket Mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__REPO_ROOT__/operator/bridges/slack
ExecStart=__NODE_BIN__ daemon.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
EnvironmentFile=-%h/.config/corvin-voice/service.env

[Install]
WantedBy=default.target
