# operator/voice/scripts/systemd/corvin-supply-chain-critical.service
[Unit]
Description=Corvin-OS daily CRITICAL-diff supply-chain alert (ADR-0021 Phase 31.3)

[Service]
Type=oneshot
ExecStart=__PYTHON_BIN__ __PLUGIN_ROOT__/scripts/supply_chain_verify.py critical --quiet
StandardOutput=journal
StandardError=journal
EnvironmentFile=-%h/.config/corvin-voice/service.env
