# operator/voice/scripts/systemd/corvin-user-style.service
[Unit]
Description=Corvin-OS daily user-style learner sweep (Layer 26)

[Service]
Type=oneshot
ExecStart=__PYTHON_BIN__ __REPO_ROOT__/operator/bridges/shared/user_style.py sweep
StandardOutput=journal
StandardError=journal
EnvironmentFile=-%h/.config/corvin-voice/service.env
