# operator/voice/scripts/systemd/corvin-audit-rotate.timer
[Unit]
Description=Daily Corvin-OS audit-rotation + retention at 03:15 local

[Timer]
OnCalendar=*-*-* 03:15:00
AccuracySec=10min
Persistent=true
Unit=corvin-audit-rotate.service

[Install]
WantedBy=timers.target
