# operator/bridges/shared/systemd/corvin-bg-monitor.service
[Unit]
Description=corvin-voice background-task session wakeup monitor
After=default.target

[Service]
Type=oneshot
# bg_monitor.py lives in .../bridges/shared, NOT .../bridges. The old
# ExecStart=__BRIDGES_DIR__/bg_monitor.py resolved to a non-existent path and
# the timer exited status=2 (ENOENT) on every tick, so the durable-completion
# backup poller never ran. Run from shared/ like the adapter unit, and share
# the adapter's service.env so CORVIN_HOME matches the writer (else the poller
# reads a different ~/.corvin tree than the one completions are queued in).
WorkingDirectory=__BRIDGES_DIR__/shared
ExecStart=/usr/bin/env python3 bg_monitor.py
StandardOutput=journal
StandardError=journal
EnvironmentFile=-%h/.config/corvin-voice/service.env
