# operator/bridges/shared/systemd/corvin-bg-monitor.timer
[Unit]
Description=Run corvin background-task wakeup monitor every minute

[Timer]
OnBootSec=3min
# OnActiveSec anchors the schedule to timer activation: without it a timer
# (re)started mid-uptime never fires again — OnBootSec lies in the past and
# OnUnitActiveSec has no anchor until the service ran at least once.
OnActiveSec=60s
OnUnitActiveSec=60s
AccuracySec=10s
Unit=corvin-bg-monitor.service
Persistent=true

[Install]
WantedBy=timers.target
