# ops/supervisord.conf
; Corvin supervisord — orchestrates adapter + bridge daemons inside one
; container. Bridges live in /opt/corvin-repo/operator/bridges/ and rely
; on shared filesystem state under /home/corvin (bind-mounted from host).
;
; Each program runs as the `corvin` user; root supervisord only does
; signal-routing.
;
; Enabling/disabling a bridge: set the matching CORVIN_BRIDGE_<name>=1|0
; env-var. Default-off so a fresh, generic container only starts the
; adapter + console; the operator opts each bridge in.

[supervisord]
nodaemon=true
user=root
logfile=/var/log/corvin/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=root:corvin

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── core: adapter (always on) ────────────────────────────────────────────
[program:adapter]
command=/opt/corvin-venv/bin/python /opt/corvin-repo/operator/bridges/shared/adapter.py
directory=/opt/corvin-repo
user=corvin
autostart=true
autorestart=true
stdout_logfile=/var/log/corvin/adapter.log
stderr_logfile=/var/log/corvin/adapter.err
environment=HOME="/home/corvin",CORVIN_HOME="/home/corvin/.corvin",PATH="/opt/corvin-venv/bin:/usr/local/bin:/usr/bin:/bin"

; ── gateway (always on) — REST API + console UI mounted under /v1/console
[program:gateway]
command=/opt/corvin-venv/bin/uvicorn corvin_gateway.app:app --host 0.0.0.0 --port 8000
directory=/opt/corvin-repo/core/gateway
user=corvin
autostart=true
autorestart=true
stdout_logfile=/var/log/corvin/gateway.log
stderr_logfile=/var/log/corvin/gateway.err
environment=HOME="/home/corvin",CORVIN_HOME="/home/corvin/.corvin",PYTHONPATH="/opt/corvin-repo/core/gateway:/opt/corvin-repo/operator/forge:/opt/corvin-repo/core/console:/opt/corvin-repo/core/compliance:/opt/corvin-repo/core/license:/opt/corvin-repo/core/compute:/opt/corvin-repo/core/delegate:/opt/corvin-repo",PATH="/opt/corvin-venv/bin:/usr/local/bin:/usr/bin:/bin",REMOTE_ORIGINS_DIR="/home/corvin/.corvin/global/a2a/origins",REMOTE_ENDPOINTS_DIR="/home/corvin/.corvin/global/a2a/endpoints"

; ── bridges (opt-in per env-var) ────────────────────────────────────────
[program:bridge-telegram]
command=node /opt/corvin-repo/operator/bridges/telegram/daemon.js
directory=/opt/corvin-repo/operator/bridges/telegram
user=corvin
autostart=%(ENV_CORVIN_BRIDGE_TELEGRAM)s
autorestart=unexpected
exitcodes=0
stdout_logfile=/var/log/corvin/bridge-telegram.log
stderr_logfile=/var/log/corvin/bridge-telegram.err
environment=HOME="/home/corvin",CORVIN_HOME="/home/corvin/.corvin",NODE_ENV="production"

[program:bridge-discord]
command=node /opt/corvin-repo/operator/bridges/discord/daemon.js
directory=/opt/corvin-repo/operator/bridges/discord
user=corvin
autostart=%(ENV_CORVIN_BRIDGE_DISCORD)s
autorestart=unexpected
exitcodes=0
stdout_logfile=/var/log/corvin/bridge-discord.log
stderr_logfile=/var/log/corvin/bridge-discord.err
environment=HOME="/home/corvin",CORVIN_HOME="/home/corvin/.corvin",NODE_ENV="production"

[program:bridge-whatsapp]
command=node /opt/corvin-repo/operator/bridges/whatsapp/daemon.js
directory=/opt/corvin-repo/operator/bridges/whatsapp
user=corvin
autostart=%(ENV_CORVIN_BRIDGE_WHATSAPP)s
autorestart=unexpected
exitcodes=0
stdout_logfile=/var/log/corvin/bridge-whatsapp.log
stderr_logfile=/var/log/corvin/bridge-whatsapp.err
environment=HOME="/home/corvin",CORVIN_HOME="/home/corvin/.corvin",NODE_ENV="production"

[program:bridge-slack]
command=node /opt/corvin-repo/operator/bridges/slack/daemon.js
directory=/opt/corvin-repo/operator/bridges/slack
user=corvin
autostart=%(ENV_CORVIN_BRIDGE_SLACK)s
autorestart=unexpected
exitcodes=0
stdout_logfile=/var/log/corvin/bridge-slack.log
stderr_logfile=/var/log/corvin/bridge-slack.err
environment=HOME="/home/corvin",CORVIN_HOME="/home/corvin/.corvin",NODE_ENV="production"

[program:bridge-email]
command=node /opt/corvin-repo/operator/bridges/email/daemon.js
directory=/opt/corvin-repo/operator/bridges/email
user=corvin
autostart=%(ENV_CORVIN_BRIDGE_EMAIL)s
autorestart=unexpected
exitcodes=0
stdout_logfile=/var/log/corvin/bridge-email.log
stderr_logfile=/var/log/corvin/bridge-email.err
environment=HOME="/home/corvin",CORVIN_HOME="/home/corvin/.corvin",NODE_ENV="production"
