# operator/bridges/whatsapp/systemd/corvin-voice-bridge-whatsapp.service
[Unit]
Description=corvin-voice WhatsApp bridge daemon (Baileys)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=__REPO_ROOT__/operator/bridges/whatsapp
ExecStart=__NODE_BIN__ daemon.js
Restart=always
RestartSec=5
TimeoutStopSec=20
StandardOutput=journal
StandardError=journal
EnvironmentFile=-%h/.config/corvin-voice/service.env

[Install]
WantedBy=default.target
