# operator/bridges/systemd/corvin-hermes-health.service
[Unit]
Description=Corvin Hermes Health Check & Repair
Documentation=https://github.com/CorvinLabs/CorvinOS/docs/claude-ref
PartOf=corvin-voice-bridge-adapter.service
After=corvin-voice-bridge-adapter.service
# StartLimit* MUST live in [Unit], not [Service] (systemd >= 229 rejects them in
# [Service] and the rate-limit is silently ignored). Bound restart storms here.
StartLimitInterval=300
StartLimitBurst=3

[Service]
Type=oneshot
# NOTE: This is a `systemctl --user` unit — it already runs as the invoking
# user. Do NOT add a `User=` / `Group=` directive: they are rejected in a
# --user unit and make systemd fail the job with status=216/GROUP (setgid
# denied), which silently disabled this healing cycle on every Linux box.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=corvin-hermes-health

# CORVIN_HOME (and any other runtime pins) come from the same service.env the
# bridge adapter reads, so reader == writer. The `-` prefix makes it optional
# so a box without the file still runs (the console falls back to its own
# CORVIN_HOME resolver).
EnvironmentFile=-%h/.config/corvin-voice/service.env
Environment="CORVIN_ACO_L5_RISKY=1"
Environment="PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# The interpreter and working-directory placeholders below are replaced at
# install time by bridge.sh install_units() with the absolute python + repo
# paths it resolved. Using the resolved interpreter (not a bare `python3`)
# means the console package is importable even on a fresh
# `uv tool install corvinos` box, where no editable/user-site install exists
# on the default python3.
WorkingDirectory=__REPO_ROOT__
ExecStart=__PYTHON_BIN__ -c "from corvin_console.aco import repair_actions; repair_actions.run_default()"

# Restart on transient failures (network hiccup, Ollama momentarily down)
# but not on permanent issues (e.g. no ollama binary at all). The restart
# rate-limit (StartLimit*) is declared in [Unit] above, as systemd requires.
Restart=on-failure
RestartSec=10
